# 7.3/bullseye/rootfs/etc/s6-overlay/s6-rc.d/laravel-cron/run
#!/command/with-contenv sh

echo "* * * * * php ${APPLICATION_ROOT:-/srv/app}/artisan schedule:run" | crontab -u www-data -

command -v crond >/dev/null && crond -f || cron -f
